# consumer/kb-serve.service
[Unit]
Description=Keyboard socket server and TCP bridge
After=network-online.target sys-kernel-config.mount
Wants=network-online.target sys-kernel-config.mount

[Service]
Type=simple
Environment=PATH=/usr/sbin:/usr/bin:/bin
ExecStartPre=/usr/sbin/gadget-bind
ExecStart=/usr/sbin/kb-serve
ExecStopPost=/usr/sbin/gadget-unbind
Restart=on-failure
RestartSec=2s

[Install]
WantedBy=multi-user.target
